# service/foot-server/run
#!/bin/sh

exec foot --server="${FOOT_SOCK}" --print-pid=3 > /dev/null

