# service/ssh-agent/finish
#!/bin/sh

rm -f "${SSH_AUTH_SOCK}"/ssh-agent.sock
exit 0
