# service/ssh-agent/run
#!/bin/sh


exec >/dev/null
exec ssh-agent -Da ${SSH_AUTH_SOCK}
